# systemd/ai-stack.service
[Unit]
Description=AI Development Stack (OpenWebUI + SearXNG)
Documentation=https://github.com/danielpsf/ai-experiments
Requires=docker.service ollama.service
After=docker.service ollama.service
BindsTo=ollama.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=daniel
Group=daniel
WorkingDirectory=/home/daniel/Documents/projects/ai-experiments
ExecStart=/home/daniel/Documents/projects/ai-experiments/scripts/management/docker-compose-wrapper.sh up
ExecStop=/home/daniel/Documents/projects/ai-experiments/scripts/management/docker-compose-wrapper.sh down
ExecReload=/home/daniel/Documents/projects/ai-experiments/scripts/management/docker-compose-wrapper.sh restart
TimeoutStartSec=300
TimeoutStopSec=60

# Environment
Environment="COMPOSE_PROJECT_NAME=ai-experiments"

[Install]
WantedBy=multi-user.target