# systemd/ollama.service
[Unit]
Description=Ollama AI Model Server
Documentation=https://ollama.ai/
After=network.target

[Service]
Type=exec
User=daniel
Group=daniel
WorkingDirectory=/home/daniel
Environment="OLLAMA_HOST=127.0.0.1:11434"
Environment="OLLAMA_ORIGINS=http://localhost:3001,http://127.0.0.1:3001"
Environment="OLLAMA_MODELS=/home/daniel/.ollama/models"
Environment="OLLAMA_KEEP_ALIVE=24h"
Environment="OLLAMA_NUM_PARALLEL=2"
Environment="OLLAMA_MAX_LOADED_MODELS=2"
Environment="OLLAMA_FLASH_ATTENTION=1"
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ollama

# Security settings (relaxed for compatibility)
NoNewPrivileges=yes
PrivateTmp=yes
ReadWritePaths=/home/daniel/.ollama

[Install]
WantedBy=multi-user.target